# etc/cron.d/print_bestall
PATH=/sbin:/usr/sbin:/bin:/usr/bin

* * * * * root /apps/bestall/print_bestall/scripts/print_bestall.sh >> /apps/bestall/logs/print_bestall.log 2>&1
